# datalocale/files/datalocale-supervisor.conf
{% from "datalocale/map.jinja" import datalocale with context -%}
[program:datalocale]
command=cubicweb-ctl start -D {{ datalocale.instance.name }}
user={{ datalocale.instance.user }}
redirect_stderr=true
stdout_file=/tmp/{{ datalocale.instance.name }}.log
environment=CW_MODE=user
numprocs=1
autostart=true
autorestart=true
startsecs=10
