# datalocale/files/saemref-supervisor.conf
{% from "datalocale/map.jinja" import datalocale with context -%}
[program:datalocale]
{% if datalocale.instance.wsgi -%}
command=uwsgi --ini /home/{{ datalocale.instance.user }}/etc/cubicweb.d/{{ datalocale.instance.name }}/uwsgi.ini
{% else -%}
command=cubicweb-ctl start -D {{ datalocale.instance.name }}
{% endif -%}
user={{ datalocale.instance.user }}
redirect_stderr=true
stdout_logfile=/var/log/cubicweb/{{ datalocale.instance.name }}.log
environment=CW_MODE=user,HOME=/home/{{ datalocale.instance.user }}
numprocs=1
autostart=true
autorestart=true
startsecs=10
